# deploy/autohome.service
[Unit]
Description=Autohome service

[Service]
Type=simple
ExecStart=/opt/autohome/start-service.sh

[Install]
WantedBy=multi-user.target
